# debian/barbican-worker.upstart
# Barbican Worker upstart script
# Used in deb build.  Keep in sync with etc/init/barbican-worker.conf

description "Barbican Key Management Workers"
start on runlevel [2345]
stop on runlevel [06]

script
    /usr/bin/barbican-worker
end script
