# debian/barbican-api.upstart
# Barbican upstart script
# Used in deb build.  Keep in sync with etc/init/barbican-api.conf

description "Barbican Key Management APIs"
start on runlevel [2345]
stop on runlevel [06]

script
    uwsgi --master --die-on-term --emperor /etc/barbican/vassals --logto /var/log/barbican/barbican-api.log --stats localhost:9314
end script
